# mkosi.packages/initrd-utils/99-add-cdrom-partitions.rules
SUBSYSTEM=="block", KERNEL=="sr[0-9]*", ACTION=="add", RUN+="/usr/sbin/kpartx -ar /dev/%k"
