# mkosi.packages/incusos-initrd-utils/initrd-boot-message.service.in
[Unit]
Description=Display early boot message
Before=boot.mount
DefaultDependencies=no

[Service]
Type=oneshot
RemainAfterExit=yes

Environment=OS_NAME="@OSNAME@"
Environment=TTYS="/dev/tty1 /dev/ttyS0"

ExecStart=/usr/bin/initrd-boot-message.sh

[Install]
WantedBy=boot.mount
