# mkosi.images/base/mkosi.extra/usr/lib/systemd/system/prometheus-node-exporter.service.d/override.conf
[Service]
Environment="ARGS=--web.listen-address=localhost:9100"
EnvironmentFile=
User=root
