# mkosi.images/incus/mkosi.extra/usr/lib/systemd/system/incus.service.d/99-disable-logfile.conf
[Service]
ExecStart=
ExecStart=/opt/incus/lib/systemd/incusd --group incus-admin $INCUS_OPTS --logfile ""
