# mkosi.images/base/mkosi.extra/usr/lib/systemd/system/kpx.service
[Unit]
Description=kpx proxy daemon

[Service]
User=nobody
Group=nogroup
ExecStart=/usr/local/bin/kpx -c /etc/kpx.yaml
Restart=on-failure
