# mkosi.images/base/mkosi.extra/usr/lib/systemd/system/apparmor.service.d/override.conf
[Unit]
After=systemd-tmpfiles-setup.service
