# mkosi.images/incus-linstor/mkosi.extra/usr/lib/systemd/system/linstor-satellite.service.d/99-dont-save-error-reports.conf
[Service]
ExecStart=
ExecStart=/usr/share/linstor-server/bin/Satellite --logs=/tmp/linstor-satellite --config-directory=/etc/linstor
