# mkosi.packages/initrd-tmpfs-root/initrd-startup-checks.service.in
[Unit]
Description=Run startup checks
After=boot.mount veritysetup-pre.target
Requires=boot.mount
Before=veritysetup.target
DefaultDependencies=no

[Service]
Type=oneshot

Environment=OS_NAME="@OSNAME@"
Environment=TTYS="/dev/tty1 /dev/ttyS0"

ExecStart=/usr/bin/initrd-startup-checks.sh

[Install]
WantedBy=veritysetup.target
