# mkosi.packages/initrd-tmpfs-root/initrd-boot-message.service.in
[Unit]
Description=Display startup message
Before=veritysetup.target
After=veritysetup-pre.target
DefaultDependencies=no

[Service]
Type=oneshot

Environment=OS_NAME="@OSNAME@"
Environment=TTYS="/dev/tty1 /dev/ttyS0"

ExecStart=/usr/bin/initrd-boot-message.sh

[Install]
WantedBy=veritysetup.target
