# mkosi.packages/initrd-tmpfs-root/initrd-message.service.in
[Unit]
Description=Display startup message
Before=veritysetup.target
After=veritysetup-pre.target
DefaultDependencies=no

[Service]
Type=oneshot

Environment=MSG="@OSNAME@ is starting..."
Environment=TTYS="/dev/tty1 /dev/ttyS0"

ExecStart=/bin/sh -c "for TTY in $TTYS; do echo $MSG > $TTY; done"

[Install]
WantedBy=veritysetup.target
